# scripts/umst-ucrs.service
[Unit]
Description=UMST Universal Calendar Resolution Spine daemon
Documentation=https://github.com/tytolabs/umst-ucrs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=umst
Group=umst
ExecStart=/usr/local/bin/umst-ucrs --peer-id 1 --port 4001 --budget-bits 20 --metrics-addr 0.0.0.0:9090
Environment=UMST_UCRS_METRICS=1
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
